# official_writeup/misc-mario/game/launcher/conf.d/ui.conf
[program:ui]
command=node index.mjs
directory=/root/launcher/ui
autorestart=true